# mcelog.cron
#!/bin/bash
/usr/sbin/mcelog --ignorenodev --filter >> /var/log/mcelog
